# debs/my-mon/etc/cron.d/backup_mon
# m h dom mon dow user command
*/5 * * * * root /opt/my-scripts/backup/backup.sh /etc/prometheus mon_conf /opt/backups
*/5 * * * * root /opt/my-scripts/backup/backup.sh /var/lib/prometheus mon_data /opt/backups
