# debs/my-mgmt/etc/cron.d/backup_mgmt
# m h dom mon dow user command
*/5 * * * * root /opt/my-scripts/backup/backup.sh /opt/my-scripts mgmt
