# debs/my-vpn/etc/cron.d/backup_vpn
# m h dom mon dow user command
*/5 * * * * root /opt/my-scripts/backup/backup.sh /etc/openvpn vpn
